# docker/etc/services.d/gpnu/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid root bash /app/run.sh